# ops/systemd/russiaww-parser-run-once.service
[Unit]
Description=RussiaWW parser-agent production run
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=parser
Group=parser
WorkingDirectory=/opt/russiaww-parser/current
Environment=APP_ENV=prod
Environment=APP_DEBUG=0
EnvironmentFile=-/opt/russiaww-parser/shared/.env.local
ExecStart=/usr/bin/php8.4 /opt/russiaww-parser/current/bin/console parser:production:run-once --limit-per-assignment=1 --no-interaction
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6

[Install]
WantedBy=multi-user.target
